# deploy/systemd/adsh-runner.service
[Unit]
Description=Run the ADSH script to check for alerts

[Service]
Type=oneshot
User=adsh-admin
Group=adsh-admin

# Set the working directory to the project root on the server
WorkingDirectory=/home/adsh-admin/adsh-project

# Define the command to execute the wrapper script
# Make sure run_adsh.sh is executable (chmod +x)
ExecStart=/home/adsh-admin/adsh-project/run_adsh.sh

# Optional: Define environment variables if needed
# Environment="PATH=/opt/adsh-project/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# Environment="YOUR_VAR=your_value"

[Install]
WantedBy=multi-user.target
